# trojan_penguin.service
[Unit]
Description = XXX

[Service]
ExecStart=/usr/bin/trojan_penguin.sh

[Install]
WantedBy=multi-user.target
